# images/cdh-cm-secondary-cdh5120/files/grafana/deploy/grafana-supervisor.conf.j2
[program:grafana]
command={{ deployDir }}/bin/grafana-server
directory={{ deployDir }}
process_name=%(program_name)s
autostart=false
autorestart=true
redirect_stderr=true
startsecs=1
