# images/cdh-cm-secondary-cdh5120/files/opentsdb/deploy/opentsdb-supervisor.conf.j2
[program:opentsdb]
command={{ deployDir }}/bin/tsdb tsd
directory={{ deployDir }}
process_name=%(program_name)s
autostart=false
autorestart=true
redirect_stderr=true
startsecs=1
